# gpu_stats.cron
*/1 * * * *	root	/usr/local/gpu_stats.sh > /tmp/gpu_stats.log 2>&1
